# etc/services.d/nginx-auth/run
#!/command/execlineb -P
# vi: ft=sh

with-contenv

/usr/local/bin/tailscale.nginx-auth -sockpath="$TS_SOCKET"
